# files/gozerbot.cron
#!/bin/sh

PID=`cat gozerbot.pid`

if test -r gozerbot.pid
then
  if `kill -CHLD $PID 1>/dev/null 2>&1`
  then
    exit 0
  else
    ./runbot 1>log 2&1 &
  fi
fi
